# .config/systemd/user/clean_journal.timer
[Unit]
Description=Clean journal directory weekly timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target