# .config/systemd/user/del_cache.timer
[Unit]
Description=Cache delete weekly timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target