# .config/systemd/user/orphan_pack.service
[Unit]
Description=Clean orphan packages

[Service]
Type=simple
ExecStart=/home/ed/bin/orphan_pack.sh
