# .config/systemd/user/test.service
[Unit]
Description=Test 

[Service]
Type=simple
ExecStart=/home/ed/bin/test.sh
